# service/camera.service
[Unit]
Description=lidar_start
[Service]
ExecStart=/usr/local/bin/aikit_shell/camera.sh

