# service/chassis.service
[Unit]
Description=chassis_start
[Service]
ExecStart=/usr/local/bin/aikit_shell/chassis_start.sh

