# service/lidar.service
[Unit]
Description=lidar_start
[Service]
ExecStart=/usr/local/bin/aikit_shell/lidar_start.sh

