# service/pcl_follow.service
[Unit]
Description=voice_start
[Service]
ExecStart=/usr/local/bin/aikit_shell/pcl_follow.sh

